# deploy/systemd/trader-btc-combine.service
[Unit]
Description=S26 Combine BTC Paper Trader (MBT 50k Eval)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Silver-Bullet-ML-BMAD
ExecStart=/root/Silver-Bullet-ML-BMAD/.venv/bin/python src/research/s26_combine_bot.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/root/Silver-Bullet-ML-BMAD/logs/s26_combine_bot.log
StandardError=append:/root/Silver-Bullet-ML-BMAD/logs/s26_combine_bot.log

[Install]
WantedBy=multi-user.target
